# crontab
* * * * * /usr/local/bin/kube-set-route.sh
